# deploy/wxparser-fixspelling.timer
[Unit]
Description=Run wxparser city-name autocorrect nightly at midnight

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true

[Install]
WantedBy=timers.target
