# deploy/wxparser-api.service
[Unit]
Description=wxparser LAN query API (reads the Postgres store)
After=network-online.target postgresql.service wxparser.service
Wants=network-online.target postgresql.service

[Service]
Type=simple
User=creed
Group=creed
WorkingDirectory=/home/creed/wxparser
Environment=HOME=/home/creed
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 -m wxparser.api
Restart=always
RestartSec=5
Nice=10

[Install]
WantedBy=multi-user.target
