# deploy/wxparser-fixterms.timer
[Unit]
Description=Run wxparser transcript term cleanup nightly (backstop; corrections also apply at write time)

[Timer]
# Staggered 30min after the city-name cleanup so the two don't run at once.
OnCalendar=*-*-* 00:30:00
Persistent=true

[Install]
WantedBy=timers.target
