# deploy/wxparser-fixspelling.timer
[Unit]
Description=Run wxparser city-name autocorrect every minute (write-time fix is primary; this mops up history + anything the slot pass misses)

[Timer]
OnCalendar=*-*-* *:*:00
AccuracySec=1s
Persistent=true

[Install]
WantedBy=timers.target
