# deploy/wxparser-agc.service
[Unit]
Description=wxparser — capture AGC (keep the input level in the decoder's sweet spot)
# After= for boot ordering only. NOT Wants=wxparser.service: this runs every 3min,
# so a Wants would resurrect a deliberately-stopped pipeline within one cycle and
# make maintenance stops impossible (`systemctl stop wxparser` would silently come
# back, then thrash against a manual arecord). The pipeline keeps itself up on its
# own (Restart=always + enabled); the AGC only reads its health.json and, when that
# heartbeat is stale (capture down), correctly makes no gain change.
After=wxparser.service

[Service]
Type=oneshot
# Runs as root: `amixer` adjusts the card and `alsactl store` writes
# /var/lib/alsa/asound.state so the gain survives a reboot. Reads the pipeline's
# world-readable health.json for the level — it never opens the audio device, so
# it does not contend with the capture service's arecord.
Environment=WX_OUT_DIR=/home/creed/wxparser/transcripts
WorkingDirectory=/home/creed/wxparser
ExecStart=/usr/bin/python3 /home/creed/wxparser/deploy/agc.py
